# .img/files/etc/init.d/fcn
#!/bin/sh /etc/rc.common
START=99
start(){
  /opt/fcn/fcn
}
stop(){
  killall fcn
}
restart(){
  killall fcn
  /opt/fcn/fcn
}
